# deploy/harness-nav.service
[Unit]
Description=Harness Navigation System
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=debian
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/debian/.Xauthority
WorkingDirectory=/home/debian/HarnessNav
ExecStart=/usr/bin/python3 /home/debian/HarnessNav/harness_nav/scripts/run_hardware.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
